# ansible/roles/monitoring/templates/redis_exporter.service.j2
[Unit]
Description=Prometheus Redis Exporter
Documentation=https://github.com/oliver006/redis_exporter
After=network.target redis.service
Requires=network.target

[Service]
Type=simple
User={{ redis_exporter_user }}
Group={{ redis_exporter_user }}
ExecStart=/usr/local/bin/redis_exporter --redis.addr=redis://localhost:6379 --web.listen-address=:{{ redis_exporter_port }}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
RestartSec=3
TimeoutStopSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/{{ redis_exporter_user }}

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target